# files/supervisor/conf.d/cron.conf
; Supervisor configuration for cron
[program:cron]
command = /usr/sbin/cron -f
