# files/supervisor/conf.d/rsyslog.conf
; supervisor configuration for rsyslog
[program:rsyslog]
command=/usr/sbin/rsyslogd -u syslog -n
